# deploy/linux/systemd/descos-backend.service
[Unit]
Description=DeskOS Backend (API + WebSocket, Port 3001)
Documentation=https://github.com/DarkWixxard/DeskOS
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=__DESCOS_USER__
Group=__DESCOS_USER__
WorkingDirectory=__DESCOS_DIR__
Environment=NODE_ENV=production
Environment=PORT=3001
Environment=PATH=__NODE_DIR__:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStart=__NPM_BIN__ run start --workspace=apps/backend
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=descos-backend

[Install]
WantedBy=multi-user.target
